# spotifypro.socket.j2
[Unit]
Description=Spotify Socket

[Socket]
ListenStream=/run/gunicorn/socket
SocketUser=www-data

[Install]
WantedBy=sockets.target